# ansible/roles/clamav/files/clamav-scan.timer
[Unit]
Description=Run ClamAV scan daily

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
